# supervisord.conf
[supervisord]
nodaemon=true

[program:npm]
command=npm start

[program:python]
command=python ../../server.py